# marianska-chata.service
[Unit]
Description=Marianska Chata Reservation System
After=network.target

[Service]
Type=forking
User=marianska
WorkingDirectory=/home/marianska/marianska
ExecStart=/home/marianska/marianska/startup-daemon.sh
ExecStop=/home/marianska/marianska/stop.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target